# sigma/rules/custom_detection_rules.yml
title: Impossible Travel for Privileged Account
id: 2d8d6b2a-61e4-4c58-9f2d-ff0d8d9d0001
status: experimental
description: Detects privileged account sign-ins from geographically distant locations within an unrealistic timeframe.
author: autobot786
date: 2026/02/08
references:
  - https://attack.mitre.org/techniques/T1078/
tags:
  - attack.t1078
  - attack.initial_access
logsource:
  category: authentication
  product: cloud
detection:
  selection:
    user.roles|contains:
      - Admin
      - Global Administrator
      - Privileged Role Administrator
  condition: selection
falsepositives:
  - Legitimate VPN/proxy egress or travel; handle with known corporate egress IP allowlists and travel confirmation.
level: high
---
title: MFA Fatigue Pattern - Multiple Failures Followed by Success
id: 2d8d6b2a-61e4-4c58-9f2d-ff0d8d9d0002
status: experimental
description: Detects repeated MFA challenge failures followed by a successful authentication.
author: autobot786
date: 2026/02/08
references:
  - https://attack.mitre.org/techniques/T1621/
tags:
  - attack.t1621
logsource:
  category: authentication
  product: cloud
detection:
  selection_fail:
    event.action: mfa_challenge
    event.outcome|contains:
      - fail
      - denied
  selection_success:
    event.action: login
    event.outcome: success
  timeframe: 15m
  condition: selection_fail and selection_success
falsepositives:
  - Users retrying MFA due to device/notification issues.
level: high
---
title: PowerShell with Base64 Encoded Command
id: 2d8d6b2a-61e4-4c58-9f2d-ff0d8d9d0003
status: experimental
description: Detects PowerShell execution with encoded commands or Base64 decode primitives.
author: autobot786
date: 2026/02/08
references:
  - https://attack.mitre.org/techniques/T1059/001/
tags:
  - attack.t1059.001
logsource:
  category: process_creation
  product: windows
detection:
  selection_img:
    process.name:
      - powershell.exe
      - pwsh.exe
  selection_cmd:
    process.command_line|contains:
      - -enc
      - -encodedcommand
      - FromBase64String
  condition: selection_img and selection_cmd
falsepositives:
  - Some admin scripts; allowlist known automation and enforce signed scripts where possible.
level: high
---
title: Service Account Interactive Login
id: 2d8d6b2a-61e4-4c58-9f2d-ff0d8d9d0004
status: experimental
description: Detects service accounts performing interactive logons.
author: autobot786
date: 2026/02/08
references:
  - https://attack.mitre.org/techniques/T1078/003/
tags:
  - attack.t1078.003
logsource:
  category: authentication
  product: windows
detection:
  selection:
    user.name|re: "^svc_|^service_|_svc$"
    logon.type:
      - Interactive
      - RemoteInteractive
  condition: selection
falsepositives:
  - Misconfigured services; deny interactive logon for service accounts.
level: high
---
title: Privileged Role Assignment Outside Change Window
id: 2d8d6b2a-61e4-4c58-9f2d-ff0d8d9d0005
status: experimental
description: Detects privileged role assignments (cloud/IAM) that should be restricted to approved change windows.
author: autobot786
date: 2026/02/08
references:
  - https://attack.mitre.org/techniques/T1098/
tags:
  - attack.t1098
logsource:
  category: iam
  product: cloud
detection:
  selection:
    event.action|contains:
      - RoleAssignmentCreated
      - Add member to role
      - AddUserToGroup
    group.name|contains:
      - Admins
      - Privileged
  condition: selection
falsepositives:
  - Emergency changes; correlate with ticketing/change management.
level: high
---
title: Scheduled Task Created with Suspicious Flags
id: 2d8d6b2a-61e4-4c58-9f2d-ff0d8d9d0006
status: experimental
description: Detects suspicious scheduled task creation (SYSTEM/highest privileges or Windows task paths).
author: autobot786
date: 2026/02/08
references:
  - https://attack.mitre.org/techniques/T1053/
tags:
  - attack.t1053
logsource:
  category: process_creation
  product: windows
detection:
  selection_img:
    process.name: schtasks.exe
  selection_cmd:
    process.command_line|contains:
      - /create
      - /sc
  suspicious:
    process.command_line|contains:
      - "/ru SYSTEM"
      - "/rl HIGHEST"
      - "\Microsoft\Windows\"
  condition: selection_img and selection_cmd and suspicious
falsepositives:
  - Legitimate scheduled jobs; baseline known task names and authors.
level: medium
---
title: Ransomware Signal - Mass Deletion Followed by File Extension Changes
id: 2d8d6b2a-61e4-4c58-9f2d-ff0d8d9d0007
status: experimental
description: Detects mass file deletions followed by rename activity to suspicious extensions.
author: autobot786
date: 2026/02/08
references:
  - https://attack.mitre.org/techniques/T1486/
tags:
  - attack.t1486
logsource:
  category: file_event
  product: windows
detection:
  deletion:
    file.action:
      - delete
      - remove
  rename:
    file.action:
      - rename
      - move
  suspicious_ext:
    file.target_extension|contains:
      - .locked
      - .encrypted
      - .enc
      - .crypt
  timeframe: 10m
  condition: deletion and rename and suspicious_ext
falsepositives:
  - Bulk file operations/migrations; validate by process, user, and host context.
level: high
---
title: Data Exfiltration - Spike to Cloud Storage Providers
id: 2d8d6b2a-61e4-4c58-9f2d-ff0d8d9d0008
status: experimental
description: Detects unusual spikes of uploads to common cloud storage providers (proxy/network logs).
author: autobot786
date: 2026/02/08
references:
  - https://attack.mitre.org/techniques/T1567/
tags:
  - attack.t1567
logsource:
  category: proxy
  product: network
detection:
  selection:
    url.domain|contains:
      - drive.google.com
      - dropbox.com
      - box.com
      - onedrive.live.com
      - mega.nz
  condition: selection
falsepositives:
  - Legitimate business usage; tune by department and device posture.
level: medium
---
title: Security Tool Tampering - AV/EDR/Firewall Disabled
id: 2d8d6b2a-61e4-4c58-9f2d-ff0d8d9d0009
status: experimental
description: Detects stopping/disabling of common security tooling services.
author: autobot786
date: 2026/02/08
references:
  - https://attack.mitre.org/techniques/T1562/001/
tags:
  - attack.t1562.001
logsource:
  category: service
  product: windows
detection:
  selection:
    service.action:
      - stop
      - disable
    service.name|contains:
      - WinDefend
      - Sense
      - WdNisSvc
      - Sysmon
      - CrowdStrike
      - CarbonBlack
  condition: selection
falsepositives:
  - Planned maintenance; correlate with change tickets.
level: high
---
title: SMB Lateral Movement Burst After Credential Use
id: 2d8d6b2a-61e4-4c58-9f2d-ff0d8d9d0010
status: experimental
description: Detects SMB access to multiple hosts in a short window (possible lateral movement).
author: autobot786
date: 2026/02/08
references:
  - https://attack.mitre.org/techniques/T1021/002/
tags:
  - attack.t1021.002
logsource:
  category: network_connection
  product: windows
detection:
  selection:
    network.protocol: smb
    event.outcome: success
  condition: selection
falsepositives:
  - Admin tools/patching/software deployment; baseline jump hosts.
level: medium
---
title: New Local or Domain Admin Created and Privileged Quickly
id: 2d8d6b2a-61e4-4c58-9f2d-ff0d8d9d0011
status: experimental
description: Detects new user creation followed by rapid addition to admin groups.
author: autobot786
date: 2026/02/08
references:
  - https://attack.mitre.org/techniques/T1136/
tags:
  - attack.t1136
logsource:
  category: iam
  product: windows
detection:
  create_user:
    event.action|contains:
      - user_created
      - NewUser
  add_admin:
    event.action|contains:
      - add_to_group
      - AddMember
    group.name|contains:
      - Domain Admins
      - Administrators
  timeframe: 30m
  condition: create_user and add_admin
falsepositives:
  - Legit onboarding; correlate with HR/ticketing.
level: high
---
title: DNS Queries to Newly Registered Domains
id: 2d8d6b2a-61e4-4c58-9f2d-ff0d8d9d0012
status: experimental
description: Detects DNS queries to newly registered domains (requires enrichment field domain_age_days).
author: autobot786
date: 2026/02/08
references:
  - https://attack.mitre.org/techniques/T1071/004/
tags:
  - attack.t1071.004
logsource:
  category: dns
  product: network
detection:
  selection:
    dns.question.name|exists: true
    threat.enrichment.domain_age_days|lte: 30
  condition: selection
falsepositives:
  - Legitimate new vendor domains; maintain allowlist.
level: medium
---
title: Suspicious Mailbox Forwarding Rule to External Recipient
id: 2d8d6b2a-61e4-4c58-9f2d-ff0d8d9d0013
status: experimental
description: Detects creation/modification of mailbox rules forwarding to external recipients.
author: autobot786
date: 2026/02/08
references:
  - https://attack.mitre.org/techniques/T1114/
tags:
  - attack.t1114
logsource:
  category: email
  product: cloud
detection:
  selection:
    event.action|contains:
      - New-InboxRule
      - Set-InboxRule
      - CreateForwardingRule
    email.forward_to|re: ".*@.*"
  filter_internal:
    email.forward_to|endswith:
      - "<YOUR_COMPANY_DOMAIN>"
  condition: selection and not filter_internal
falsepositives:
  - Approved forwarding; restrict tenant-wide and require approvals.
level: high
---
title: VPN Login from Non-Compliant Device
id: 2d8d6b2a-61e4-4c58-9f2d-ff0d8d9d0014
status: experimental
description: Detects VPN access from devices failing compliance/posture checks.
author: autobot786
date: 2026/02/08
references:
  - https://attack.mitre.org/techniques/T1133/
tags:
  - attack.t1133
logsource:
  category: vpn
  product: network
detection:
  selection:
    event.action: vpn_login
    event.outcome: success
    device.compliant: false
  condition: selection
falsepositives:
  - BYOD exceptions; enforce conditional access and compliance policies.
level: medium
---
title: Dormant User Access Spike to Sensitive Data
id: 2d8d6b2a-61e4-4c58-9f2d-ff0d8d9d0015
status: experimental
description: Detects dormant accounts accessing sensitive resources (requires enrichment last_seen_days and classification).
author: autobot786
date: 2026/02/08
references:
  - https://attack.mitre.org/techniques/T1087/
tags:
  - attack.t1087
logsource:
  category: file_access
  product: multiple
detection:
  selection:
    user.last_seen_days|gte: 30
    data.classification|contains:
      - Sensitive
      - Confidential
      - PII
  condition: selection
falsepositives:
  - Returning employees/contractors; confirm role and access need.
level: high
